# root/rpi4b/etc/systemd/system/systemd-nspawn@hass.service.d/network.conf
[Unit]
Wants=network-online.target systemd-resolved.service
After=network-online.target systemd-resolved.service
