# ssh/.config/systemd/user/authorized-keys.timer
[Unit]
Description=Download .ssh/authorized_keys

[Timer]
OnBootSec=5min
OnCalendar=*-*-* *:00:00
Persistent=true

[Install]
WantedBy=default.target
