# root/krekenava2/etc/systemd/system/mqtt.socket
[Socket]
ListenStream=1883

[Install]
WantedBy=sockets.target
